# com.claude-sessions.scheduler.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-sessions.scheduler</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{USER_HOME}}/.local/bin/claude-sessions</string>
        <string>start-now</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{USER_HOME}}/.local/share/claude-sessions/scheduler.out</string>
    <key>StandardErrorPath</key>
    <string>{{USER_HOME}}/.local/share/claude-sessions/scheduler.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{VENV_PATH}}/bin:/opt/homebrew/bin:/opt/homebrew/opt/node@22/bin:/usr/local/bin:/usr/bin:/bin:{{USER_HOME}}/.local/bin</string>
        <key>HOME</key>
        <string>{{USER_HOME}}</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>{{USER_HOME}}</string>
</dict>
</plist>